# deb/com2kids/DEBIAN/preinst
#! /bin/bash
if [ -f /etc/init.d/com2kids ]
then
	service com2kids stop
fi
update-rc.d -f com2kids remove
rm -rf /var/log/com2kids/
